# packaging/aqueous-inputd.socket
[Unit]
Description=Aqueous input daemon socket

[Socket]
ListenStream=%t/aqueous-inputd.sock
SocketMode=0600

[Install]
WantedBy=sockets.target
